# fpga/arm/RemoteSystemsTempFiles/172.20.180.95/etc/rc.local
#!/bin/bash
echo "RUN FISH COUNT"
chmod +x /home/root/appl
/home/root/./appl
exit 0