# fpga/arm/RemoteSystemsTempFiles/172.20.180.81/etc/rc.local
#!/bin/bash 
#chmod +x /home/root/run_appl.sh
/home/root/run_appl.sh
exit 0